# s6/shadowsocks/finish
#!/bin/sh

kill `ps | awk '/[s]s-server/{print $1}'`
