# s6/kcptun/finish
#!/bin/sh

kill `ps | awk '/[k]cptun/{print $1}'`
