# s6/kcptun/run
#!/bin/execlineb -P

/usr/bin/kcptun -t 127.0.0.1:8388 -l :29900 -mode fast2
